# deployment/myroom-sensor-monitor.timer.template
[Unit]
Description=Run MyRoom sensor monitor every 10 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
